# rpi-fan.service
[Unit]
Description=Raspi Fan service.

[Service]
Type=simple
ExecStart=/usr/bin/rpi-fan

[Install]
WantedBy=multi-user.target
